# etc/service/phpfpm/run
#!/usr/bin/env bash

source /utils/vault-env.sh

php-fpm${PHP_VERSION} -F -O
