# etc/service/phpfpm/run
#!/usr/bin/env bash

echo 'Starting phpfpm'

. /utils/cache-config.sh

php-fpm${PHP_VERSION} -F -O
