# etc/service/websocket/run
#!/usr/bin/env bash

echo 'Starting Websockets'

. /utils/cache-config.sh

php ${WORK_DIR}/artisan websockets:serve
