# etc/service/horizon/run
#!/usr/bin/env bash

source /utils/vault-env.sh

php ${WORK_DIR}/artisan horizon
